# utils/ssh-tunnel/ssh-tunnel@.service
[Unit]
Description=SSH Tunnel for instance %i
After=network.target

[Service]
Type=simple
#User=tunnel
WorkingDirectory=/volume1/homes/tunnel
Environment=DEBUG=true
EnvironmentFile=/volume1/homes/tunnel/tunnel.env
ExecStart=/bin/su -s /bin/sh -c '/volume1/homes/tunnel/tunnel.sh ${%i_TUNNEL} ${%i_CONNECTION}' tunnel
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
